# packaging/systemd/voicekb.service
[Unit]
Description=VoiceKB — local voice-to-text dictation daemon
Documentation=file://@VOICEKB_DIR@/README.md
After=graphical-session.target pipewire.service
Wants=graphical-session.target
StartLimitBurst=5
StartLimitIntervalSec=30s

[Service]
Type=notify
NotifyAccess=main
WorkingDirectory=@VOICEKB_DIR@
Environment=HF_HOME=@VOICEKB_DIR@/models/hf
ExecStart=@VOICEKB_DIR@/.venv/bin/python -m voicekb run-daemon --hotkey capslock
RuntimeDirectory=voicekb
Restart=on-failure
RestartSec=2s
TimeoutStartSec=60s
TimeoutStopSec=5s

# Hardening — values vetted per systemd.exec(5).
# @VOICEKB_DIR@ is substituted with the repo path by `just install-user`.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@VOICEKB_DIR@/models %t/voicekb
PrivateTmp=yes
PrivateDevices=no
RestrictAddressFamilies=AF_UNIX AF_NETLINK
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
CapabilityBoundingSet=
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=default.target
